# docker/wwwroot/crontab

5 *    * * *   zoretad /var/www/htdocs/rozinante/mail2web
